# setup/crontabs/root.template
@reboot _TARGET_/networkmanager/killswitch
